# debian/pvc-daemon.prerm
#!/bin/sh

# Disable the services
systemctl disable pvc-flush.service
systemctl disable pvcd.service
systemctl disable pvcd.target
