# debian/pvc-daemon.prerm
#!/bin/sh

# Disable the servive
systemctl disable pvcd.service

