# debian/pvc-virtualization-daemon.prerm
#!/bin/sh

# Disable the service
systemctl disable pvcvd.service

