# debian/pvc-network-daemon.prerm
#!/bin/sh

# Disable the service
systemctl disable pvcnd.service

