# debian/pvc-daemon.prerm
#!/bin/sh

# Disable the service
systemctl disable pvcd.service
