# debian/pvc-daemon-health.prerm
#!/bin/sh

# Disable the services
systemctl disable pvchealthd.service
