# debian/pvc-virtualization-daemon.postinst
#!/bin/sh

# Enable the servive
systemctl enable /lib/systemd/system/pvcd.service

echo "The PVC daemon has not been started. Create a config file at /etc/pvc/pvcd.conf then start it."
