# debian/pvc-daemon-worker.prerm
#!/bin/sh

# Disable the services
systemctl disable pvcworkerd.service
